# loader/systemd/KlipperScreen.service.d/override.conf
[Unit]
After=plymouth-quit.service
Wants=plymouth-quit.service

[Service]
ExecStartPre=/bin/sh -lc 'plymouth quit || true'
